# web/llama-vision.service
[Unit]
Description=llama-server vision sidecar (severian-vision multimodal)
After=network.target

[Service]
Type=simple
User=tmancino
Group=tmancino
WorkingDirectory=/opt/llama.cpp-vision
Environment=LD_LIBRARY_PATH=/opt/llama.cpp-vision/build/bin
ExecStart=/opt/llama.cpp-vision/build/bin/llama-server \
  --model /opt/severian/model/severian-vision-q4_k_m.gguf \
  --mmproj /opt/severian/model/mmproj-severian-vision.gguf \
  --host 127.0.0.1 --port 8081 \
  --n-gpu-layers 99 \
  --ctx-size 4096 \
  --no-warmup
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/llama-vision.log
StandardError=append:/var/log/llama-vision.log

[Install]
WantedBy=multi-user.target
